# tradbot.service
[Unit]
Description=Tradbot - Automated Options Trading System
After=network-online.target
Wants=network-online.target

# Restart on failure
StartLimitInterval=300
StartLimitBurst=5

[Service]
Type=simple

# User and group  
User=tradbot
Group=tradbot

# Working directory
WorkingDirectory=/opt/tradbot
# Load environment from external file (managed by sysadmin)
EnvironmentFile=/etc/tradbot.env

# Execution
ExecStart=/opt/tradbot/venv/bin/python /opt/tradbot/main.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradbot

# Resource limits
MemoryLimit=1G
CPUQuota=50%
WatchdogSec=60
TimeoutStopSec=30

# Restart policy
Restart=always
RestartSec=30s

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/tradbot /var/lib/tradbot

[Install]
WantedBy=multi-user.target
